# files/yandex_ddns.service
[Unit]
Description=Update Yandex DNS record

[Service]
Type=oneshot
ExecStart=/opt/yandex_ddns/update
